# infra/systemd/eww-status.service
[Unit]
Description=EWW status
[Service]
Type=oneshot
ExecStart=/git/eww/bin/eww-status
