# docs/infra/infra/systemd/eww-autocommit.service
[Unit]
Description=EWW auto-commit
[Service]
Type=oneshot
ExecStart=/git/eww/bin/eww-commit
