# core/systemd/graph-snapshot.service
[Unit]
Description=Daily Graph Snapshot Service
Documentation=file:///home/jakubc/git/eww/docs/software/obsidian/plugins/core/graph-view/graph-view-plan.md

[Service]
Type=oneshot
ExecStart=/home/jakubc/git/eww/dev/scripts/graph-snapshot.sh
WorkingDirectory=/home/jakubc/git/eww
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
